# etc/udev/rules.d/90-elphel-disk.rules
KERNEL=="sd*", RUN+="/usr/local/bin/elpheldisk"

