# etc/udev/rules.d/90-eyesis-disk.rules
KERNEL=="sd*", RUN+="/usr/local/bin/eyesisdisk"

